# ssl-cert-system-go/dist/ssl-cert-system-go-linux-1.0.2/ssl-cert-system.service
[Unit]
Description=SSL Certificate Management System
After=network.target mysql.service

[Service]
Type=simple
User=sslapp
WorkingDirectory=/opt/ssl-cert-system
ExecStart=/opt/ssl-cert-system/ssl-cert-system
Restart=always
RestartSec=5
Environment=NODE_ENV=production

[Install]
WantedBy=multi-user.target
