# deployment/dist/ssl-cert-system-v1.0.0/docker/supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid
childlogdir=/var/log/supervisor

[unix_http_server]
file=/var/run/supervisor.sock
chmod=0700

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

# Nginx服务
[program:nginx]
command=nginx -g "daemon off;"
autostart=true
autorestart=true
stderr_logfile=/var/log/sslapp/nginx_error.log
stdout_logfile=/var/log/sslapp/nginx_access.log
user=root

# SSL证书管理系统后端
[program:sslapp-backend]
command=node /app/backend/dist/simple-app.js
directory=/app/backend
autostart=true
autorestart=true
stderr_logfile=/var/log/sslapp/backend_error.log
stdout_logfile=/var/log/sslapp/backend_access.log
user=sslapp
environment=NODE_ENV="production",PORT="3001"

# 日志轮转
[program:logrotate]
command=/usr/sbin/logrotate -f /etc/logrotate.conf
autostart=false
autorestart=false
startsecs=0
